# etc/launchd/daemons/org.openzfsonosx.zed.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.openzfsonosx.zed</string>
	<key>ProgramArguments</key>
	<array>
		<string>@sbindir@/zed</string>
		<string>-vfF</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/private/var/log/org.openzfsonosx.zed.err</string>
	<key>StandardOutPath</key>
	<string>/private/var/log/org.openzfsonosx.zed.log</string>
</dict>
</plist>
