# udev/rules.d/90-zfs.rules.in
SUBSYSTEM!="block|misc", GOTO="zfs_end"
ACTION!="add|change", GOTO="zfs_end"

ENV{ID_FS_TYPE}=="zfs", RUN+="/sbin/modprobe zfs"
ENV{ID_FS_TYPE}=="zfs_member", RUN+="/sbin/modprobe zfs"

KERNEL=="null", SYMLINK+="root"
SYMLINK=="null", SYMLINK+="root"

KERNEL=="zfs", MODE="0666", OPTIONS+="static_node=zfs"

LABEL="zfs_end"
