# debian/zfs-zed.prerm
#!/bin/bash -eux

case $1 in
upgrade | remove)
	systemctl disable zfs-zed.service
	;;
esac

exit 0
