# debian/zfsutils-linux.postinst
#!/bin/bash -eux

case $1 in
configure)
	systemctl enable zfs-import-cache.service
	systemctl enable zfs-import-scan.service
	systemctl enable zfs-import.target
	systemctl enable zfs-mount.service
	systemctl enable zfs-share.service
	systemctl enable zfs.target
	;;
esac

exit 0
