# etc/systemd/system/zed.service.in
[Unit]
Description=ZFS Event Daemon (zed)
Documentation=man:zed(8)
After=zfs-import-cache.service
After=zfs-import-scan.service

[Service]
Type=forking
ExecStart=@sbindir@/zed
PIDFile=@runstatedir@/zed.pid
User=root
Group=root
Restart=on-abort
