# etc/systemd/system/zfs-share-cleanup.service.in
[Unit]
Description=ZFS file system shares cleanup
After=nfs-server.service
PartOf=nfs-server.service nfs-kernel-server.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=-/usr/lib/zfs-linux/share-cleanup

[Install]
WantedBy=zfs.target
