# etc/systemd/system/zfs-zed.service.in
[Unit]
Description=ZFS Event Daemon (zed)
Documentation=man:zed(8)
ConditionPathIsDirectory=/sys/module/zfs

[Service]
ExecStart=/usr/sbin/zed -F
Restart=on-abort

[Install]
Alias=zed.service
WantedBy=zfs.target
