# etc/systemd/system/zfs-volume-wait.service.in
[Unit]
Description=Wait for ZFS Volume (zvol) links in /dev
DefaultDependencies=no
After=systemd-udev-settle.service
After=zfs-import.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@bindir@/zvol_wait

[Install]
WantedBy=zfs-volumes.target
