# etc/systemd/system/zfs-mount.service.in
[Unit]
Description=Mount ZFS filesystems
DefaultDependencies=no
After=systemd-udev-settle.service
After=zfs-import.target
After=systemd-remount-fs.service
Before=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@sbindir@/zfs mount -a

[Install]
WantedBy=zfs-share.service
WantedBy=zfs.target
