# etc/systemd/system/zfs-module.service.in
[Unit]
Description=ZFS file system kernel modules
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe zfs
