# etc/systemd/system/zfs-trim-monthly@.timer.in
[Unit]
Description=Monthly zpool trim timer for %i
Documentation=man:zpool-trim(8)

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=1h
Unit=zfs-trim@%i.service

[Install]
WantedBy=timers.target
