# etc/systemd/system/zfs-import.target.in
[Unit]
Description=ZFS pool import target
Before=dracut-mount.service

[Install]
WantedBy=zfs.target
