# etc/systemd/system/zfs-import.target.in
[Unit]
Description=ZFS pool import target

[Install]
WantedBy=zfs-mount.service
WantedBy=zfs.target
