# etc/systemd/system/zfs-zed.service.in
[Unit]
Description=ZFS Event Daemon (zed)
Documentation=man:zed(8)
ConditionPathIsDirectory=/sys/module/zfs

[Service]
# DLPX-75168 -- to avoid zed crashes we run with jobs set to 1
ExecStart=/usr/sbin/zed -F -j 1
Restart=on-abort

[Install]
Alias=zed.service
WantedBy=zfs.target
