# etc/systemd/system/zfs-volumes.target.in
[Unit]
Description=ZFS volumes are ready
After=zfs-volume-wait.service
Requires=zfs-volume-wait.service

[Install]
WantedBy=zfs.target
