# etc/systemd/system/zfs.target.in
[Unit]
Description=ZFS startup target

[Install]
WantedBy=multi-user.target
