# etc/systemd/system/zpool@.path.in
[Unit]
Description=zpool ready path
DefaultDependencies=no

[Path]
PathExists=/dev/zpool/%i/ready
Unit=zpool@%i.service
