# contrib/dracut/90zfs/zfs-snapkeep-bootfs.service.in
[Unit]
Description=Purge old snapshots of the bootfs
Requisite=zfs-import.target
After=zfs-import.target dracut-pre-mount.service zfs-snapshot-bootfs.service zfs-rollback-bootfs.service
Before=dracut-mount.service
DefaultDependencies=no
ConditionKernelCommandLine=bootfs.snapkeep

[Service]
Type=oneshot
ExecStart=-/bin/sh -c '. /lib/dracut-zfs-lib.sh; decode_root_args || exit; [ "$root" = "zfs:AUTO" ] && root="$BOOTFS"; SNAPKEEP="$(getarg bootfs.snapkeep)"; [ "$SNAPKEEP" -ge 0 ] || exit 1; printf -v IFS "\n"; for SNAPNAME in $(/sbin/zfs list -H -o name -t snapshot "$root" | head -n -$SNAPKEEP); do echo "$SNAPNAME" | grep -q "@" || continue; /sbin/zfs destroy -Rf "$SNAPNAME"; done'
RemainAfterExit=yes
